# service/dnsmasq/run
#!/bin/sh

dnsmasq --conf-file=/etc/service/dnsmasq/dnsmasq.conf
